# debian/univention-nqrustbackup.postinst.in
#!/bin/sh
# source debconf stuff
. /usr/share/debconf/confmodule
. /usr/share/univention-lib/all.sh

exec 3>/dev/null

# ... rest of your code ...
if [ "$1" = "configure" ]; then

    # on upgrade, set backup_myself to yes if it was not set
    filestorage=$(ucr get nqrustbackup/filestorage)
    backup_myself=$(ucr get nqrustbackup/backup_myself)
    if [ -n "$filestorage" -a -z "$backup_myself" ]; then
        ucr set nqrustbackup/backup_myself=yes
    fi

    # set defaults
    ucr set \
            nqrustbackup/filestorage?/var/lib/nqrustbackup/storage \
            nqrustbackup/max_full_volume_bytes?20 \
            nqrustbackup/max_full_volumes?1 \
            nqrustbackup/max_diff_volume_bytes?10 \
            nqrustbackup/max_diff_volumes?1 \
            nqrustbackup/max_incr_volume_bytes?1 \
            nqrustbackup/max_incr_volumes?1 \
            nqrustbackup/backup_myself?no \
            nqrustbackup/webui/console/user1/username?Administrator \
            nqrustbackup/webui/console/user1/password?

    # commit defaults
    ucr commit \
            /etc/nqrustbackup/nqrustbackup_console.conf
            #/etc/nqrustbackup/nqrustbackup-dir.conf \
            #/etc/nqrustbackup/nqrustbackup-sd.conf \
            #/etc/nqrustbackup/nqrustbackup-fd.conf \

    touch /etc/nqrustbackup/autogenerated/clients.include

    # activate our listener plugin
    invoke-rc.d univention-directory-listener restart

    # restart nqrustbackup
    invoke-rc.d nqrustbackup-dir restart
    invoke-rc.d nqrustbackup-sd restart
    invoke-rc.d nqrustbackup-fd restart

    [ -x "/etc/init.d/univention-firewall" ] && invoke-rc.d univention-firewall restart

    # reload apache
    invoke-rc.d apache2 reload

    call_joinscript 62univention-nqrustbackup.inst
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
